# systemd/user/greenclip.service
[Unit]
Description=Rofi-greenclip daemon

[Service]
ExecStart=/usr/bin/greenclip daemon

[Install]
WantedBy=default.target
